# deployment/launchd/com.openclaw.signal-platform.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.signal-platform</string>

    <!--
      Interpreter: the repo-local venv Python. This ensures all our pinned
      deps (prometheus-client, structlog, pandas, yfinance, pyarrow, click)
      are available without touching the system Python.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/matthew.zhao/signal-platform/.venv/bin/python</string>
        <string>-m</string>
        <string>signal_platform.service</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/matthew.zhao/signal-platform/.venv/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SIGNAL_PLATFORM_METRICS_PORT</key>
        <string>9095</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/matthew.zhao/signal-platform</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <!-- Back off to 1-minute retry after a crash loop so we don't hammer
         the machine if something is fundamentally broken. -->
    <key>ThrottleInterval</key>
    <integer>60</integer>

    <!--
      User-writable log paths. ``/var/log/`` would need sudo to create +
      write; keeping logs in ``~/signal-platform-logs/`` lets the agent
      run entirely under the user's account. Production may later mirror
      to NAS; that's a separate concern.
    -->
    <key>StandardOutPath</key>
    <string>/Users/matthew.zhao/signal-platform-logs/service.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/matthew.zhao/signal-platform-logs/service.err</string>
</dict>
</plist>
